# deploy/etc/s6-overlay/s6-rc.d/placement/run
#!/command/with-contenv sh


exec /root/.dapr/bin/placement