# deploy/etc/s6-overlay/s6-rc.d/dapr/run
#!/command/with-contenv execlineb

exec dapr run --app-id kitchenai --app-port 8001 \
    --app-protocol http \
    -components-path /components
