# deploy/etc/s6-overlay/s6-rc.d/dapr/run
#!/command/with-contenv execlineb

exec dapr run  \
    --app-id app \
    --app-protocol http \
    --scheduler-host-address "" \
    --dapr-http-port 3500