# deploy/etc/s6-overlay/s6-rc.d/kitchenai-dapr/run
#!/command/with-contenv execlineb

exec dapr run --app-port 8001 \
    --app-id kitchenai \
    --app-protocol http \
    --scheduler-host-address "" \
    --dapr-http-port 3501 -- python3 /app/kitchenai/app.py
